# honcho/rootfs/etc/s6-overlay/s6-rc.d/honcho-api/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
bashio::log.warning "Honcho API stopped."
