# honcho/rootfs/etc/s6-overlay/s6-rc.d/honcho-api/run
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Start Honcho API server (FastAPI on port 8000)
# ==============================================================================

declare CONFIG_DIR="/data/honcho"

# Load environment
if [ -f "${CONFIG_DIR}/env" ]; then
    set -a
    # shellcheck source=/dev/null
    source "${CONFIG_DIR}/env"
    set +a
fi

cd /app

bashio::log.info "Starting Honcho API server..."
# Use uvicorn directly with our health check wrapper.
# The health module imports src.main:app and registers /health on it.
exec /app/.venv/bin/python -m uvicorn health:app \
    --host 0.0.0.0 \
    --port 8000 \
    --app-dir /app
