# honcho/rootfs/etc/s6-overlay/s6-rc.d/redis/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Graceful Redis shutdown
# ==============================================================================

bashio::log.info "Stopping Redis..."
redis-cli shutdown nosave || true
