# honcho/rootfs/etc/s6-overlay/s6-rc.d/postgresql/finish
#!/command/with-contenv bashio
# shellcheck shell=bash
# ==============================================================================
# Graceful PostgreSQL shutdown
# ==============================================================================

bashio::log.info "Stopping PostgreSQL..."
su - postgres -c "/usr/lib/postgresql/17/bin/pg_ctl stop -D /data/postgresql -m fast" || true
